# etc/init.d/automountd
#!/sbin/openrc-run

command=/usr/sbin/automountd
name="autofs Mounting Daemon"

depend()
{
	need localmount
	use net logger dns
	before inetd xinetd
}

stop_post()
{
	sleep 1
}
