# etc/init.d/automount
#!/sbin/openrc-run

name="NFS Automount Client"

depend()
{
	need localmount nfsclient automountd autounmountd
	use net logger dns
	before inetd xinetd
}

start()
{
	ebegin "Starting $name"
	/usr/sbin/automount ${automount_flags}
	eend $?
}

stop()
{
	/sbin/umount -At autofs
}
