# etc/init.d/automountd
#!/sbin/openrc-run

command=/usr/sbin/automountd
command_args="-d"
name="automountd"
pidfile="/var/run/automountdaemond.pid"
supervisor=supervise-daemon
supervise_daemon_args="-1 /var/log/automountdaemond.log -2 /var/log/automountdaemond.log"

depend()
{
        need localmount
	provide automountd
        use net logger dns
        before inetd xinetd
}
