# etc/init.d/cron
#!/sbin/openrc-run
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	use clock logger
	need localmount
}

start() {
	ebegin "Starting cron daemon"
	start-stop-daemon --start --quiet --exec /usr/sbin/cron
	eend $?
}

stop() {
	ebegin "Stopping cron daemon"
	start-stop-daemon --stop --quiet --pidfile /var/run/cron.pid
	eend $?
}
