# etc/init.d/hostapd
#!/sbin/openrc-run

name="hostapd"
description="Authenticator for IEEE 802.11 networks"
command=/usr/sbin/hostapd

devarg=${RC_SVCNAME##*.}
if [ -n "$devarg" -a "$devarg" != "devmatch" ]; then
        hostdev="$devarg"
        name="$name ($hostdev)"
	conf_file="/etc/hostapd-${hostdev}.conf"
	pidfile="/var/run/hostapd-${hostdev}.pid"
else
	conf_file="/etc/${name}.conf"
	pidfile="/var/run/${name}.pid"
fi

command_args="-P ${pidfile} -B ${conf_file}"
required_files="${conf_file}"


depend()
{
	need localmount
	use net
	keyword -shutdown -jail
}

start_pre()
{
	load_kld wlan_xauth wlan_wep wlan_tkip wlan_ccmp
}
