# etc/init.d/kdc
#!/sbin/openrc-run

. /etc/rc.subr

name=kdc

depend() {
	provide kdc
	after network
	before login
	keyword -shutdown
}

start_pre()
{
	command_args="$command_args --detach"
}

start() {
	set_rcvar_obsolete kerberos5_server_enable kdc_enable
	set_rcvar_obsolete kerberos5_server kdc_program
	set_rcvar_obsolete kerberos5_server_flags kdc_flags
}

stop() {
	return 0
}
