# etc/init.d/keyserv
#!/sbin/openrc-run

name="keyserv"
command="/usr/sbin/${name}"

depend() {
	provide keyserv
	need ypset rpcbind
	keyword -shutdown
}
