# etc/init.d/keyserv
#!/sbin/openrc-run

. /etc/rc.subr

name="keyserv"
command="/usr/sbin/${name}"

depend() {
	provide keyserv
	need ypset
	keyword -shutdown
}

star_pre()
{
	force_depend rpcbind || return 1
}
