# etc/init.d/kfd
#!/sbin/openrc-run

. /etc/rc.subr

name=kfd

depend()
{
	provide kfd
	after network
	keyword -shutdown
}

start()
{

	command_args="$command_args -i &"
}

stop()
{
	return 0
}
