# etc/init.d/ldapd
#!/sbin/openrc-run
# PROVIDE: ldapd
# REQUIRE: DAEMON FILESYSTEMS devfs
# KEYWORD: nojail shutdown

name="ldapd"
rcvar="ldapd_enable"
command="/usr/sbin/${name}"

depend() {
	need devfs
	keyword -jail -shutdown
}
