# etc/init.d/nfscbd
#!/sbin/openrc-run
#
# $FreeBSD$
#

# PROVIDE: nfscbd
# REQUIRE: NETWORKING nfsuserd
# KEYWORD: nojail shutdown

name="nfscbd"
command="/usr/sbin/${name}"

depend() {
	need nfsuserd
	keyword -jail -shutdown
}
