# etc/init.d/ypbind
#!/sbin/openrc-run

name="ypbind"
description="NIS domain binding daemon"

command="/usr/sbin/$name"
command_args="${nis_client_flags}"

depend()
{
	need localmount net
	keyword -shutdown -stop
}

start_pre()
{
	local _domain

	_domain=`domainname`
	if [ -z "$_domain" ]; then
		ewarn "NIS domainname(1) is not set."
		return 1
	fi
}
