# sbin/init/init.d/bootparams
#!/sbin/openrc-run
#
# $FreeBSD$
#

depend()
{
	provide bootparams
	need rpcbind
	keyword -jail
}

name="bootparamd"
desc="Boot parameter daemon"
required_files="/etc/bootparams"
command="/usr/sbin/${name}"
