# sbin/init/init.d/kfd
#!/sbin/openrc-run

name=kfd

depend()
{
	provide kfd
	after network
	keyword -shutdown
}

start()
{
	return 0
}

stop()
{
	return 0
}
